# .config/systemd/user/podman-auto-update.timer.d/override.conf
# Automatically update the containers once per week
[Timer]
OnCalendar=
OnCalendar=weekly
